# appliance/home/pi/.config/autostart/graphical-startup.desktop
[Desktop Entry]
Name=Graphical startup tasks
Exec=/home/pi/graphical-startup.sh
Terminal=false
Type=Application
